# P06-gateway-full/assetarc-gateway/systemd/assetarc-gateway.service

[Unit]
Description=AssetArc Gateway (Gunicorn)
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/opt/assetarc-gateway
EnvironmentFile=/opt/assetarc-gateway/.env
ExecStart=/opt/assetarc-gateway/.venv/bin/gunicorn -c gunicorn.conf.py app:app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
